# minio/rootfs/etc/services.d/minio/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the MinIO service
# ==============================================================================

# Get config values
data_path=$(bashio::config 'data_path')


# Optional SSL configuration
if bashio::config.has_value 'ssl_certificate'; then
    ssl_cert=$(bashio::config 'ssl_certificate')
    ssl_key=$(bashio::config 'ssl_key')
    export MINIO_CERT_FILE="${ssl_cert}"
    export MINIO_KEY_FILE="${ssl_key}"
fi

# Set environment variables
export MINIO_ROOT_USER="$(bashio::config 'username')"
export MINIO_ROOT_PASSWORD="$(bashio::config 'password')"


# Start MinIO
exec minio server "${data_path}" --console-address ":9001"