# minio/rootfs/etc/services.d/minio/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the MinIO service
# ==============================================================================

# Get config values
data_path=$(bashio::config 'data_path')


# SSL configuration
bashio::config.require.ssl

certfile=$(bashio::config 'certfile')
export MINIO_CERT_FILE="${certfile}"
keyfile=$(bashio::config 'keyfile')
export MINIO_KEY_FILE="${keyfile}"


# Set environment variables
export MINIO_ROOT_USER="$(bashio::config 'username')"
export MINIO_ROOT_PASSWORD="$(bashio::config 'password')"


# Start MinIO
exec minio server "${data_path}" --console-address ":9001"